# container.target
[Unit]
Description=Container target  - special target for running systemd in docker containers
